# config/includes.chroot/etc/cron.daily/update-flashplugin-nonfree
#!/bin/sh
test -x /usr/sbin/update-flashplugin-nonfree || exit 0
/usr/sbin/update-flashplugin-nonfree --install --quiet >/dev/null 2>&1
