# config/includes.chroot/etc/skel/.config/autostart/xpad.desktop
[Desktop Entry]
Type=Application
Exec=xpad
Name[fr_FR]=XPad
Name=XPad
Comment[fr_FR]=Bloc-Notes XPad
Comment=Bloc-Notes XPad
Icon=xpad
